# buildroot/board/fos-next/rootfs_overlay/etc/inittab
# /etc/inittab - BusyBox init configuration for fos-next
#
# Action descriptions:
#   sysinit   - run once at boot before anything else
#   respawn   - restart if process exits
#   askfirst  - like respawn, but ask before starting a new shell
#   once      - run once and do not restart
#   wait      - run once and wait for completion

# Mount important pseudo-filesystems
null::sysinit:/bin/mount -t proc proc /proc
null::sysinit:/bin/mount -t sysfs sysfs /sys
null::sysinit:/bin/mount -t devtmpfs devtmpfs /dev
null::sysinit:/bin/mkdir -p /dev/pts
null::sysinit:/bin/mount -t devpts devpts /dev/pts

# Run init.d scripts in order
# 'console' tells BusyBox init to open /dev/console and dup2 it to
# stdin/stdout/stderr before exec-ing rcS, so the entire process tree
# (rcS → S99fos → fos-agent) inherits a usable console fd.
console::sysinit:/etc/init.d/rcS

# fos-agent runs as PID 1 init (started by rcS/S99fos)
# Console on tty1 for debug mode (respawn only if agent exits 0)
tty1::respawn:/sbin/getty -L tty1 0 vt100
ttyS0::respawn:/sbin/getty -L ttyS0 115200 vt100

# Reboot/shutdown signals
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
